# initscripts/funtoo/pyFC
#!/sbin/runscript 

exec='python /usr/share/pyFanController/'
pidfile='/var/run/pyFC.pid'
config_file='/etc/pyFC.ini'

depend() { 
  need localmount
  after bootmisc
} 

start() { 
 ebegin "Starting fan controller" 
 start-stop-daemon --start --name pyFC --env PYFC_CONFIG_PATH="$config_file"--pidfile "$pidfile" -w 10 --background --exec $exec
 eend $? 
} 

stop() { 
 ebegin "Stopping fan controller" 
 start-stop-daemon --stop --pidfile "$pidfile"
# kill SIGTERM  
 eend $? 
}
